# tydom2mqtt_hassio/rootfs/etc/services.d/tydom2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: tydom2MQTT
# Runs tydom2MQTT
# ==============================================================================
echo Hello world!

ls /
ls /usr/bin
ls /opt/

exec python /opt/tydom2MQTT/main.py

# python3 -u main.py