# tydom2mqtt_hassio/rootfs/etc/services.d/tydom2mqtt/run
#!/usr/bin/execlineb -P
# ==============================================================================
# Home Assistant Add-on: tydom2MQTT
# Runs tydom2MQTT
# ==============================================================================

# bashio::log.info 'Starting tydom2MQTT...'

# ls
# ls /opt/
python3 -u /opt/tydom2MQTT/main.py

# python3 -u main.py