# tydom2mqtt_hassio/rootfs/etc/services.d/tydom2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: tydom2MQTT
# Runs tydom2MQTT
# ==============================================================================

ls /
ls /usr/bin
ls /opt/

/usr/bin/python /opt/tydom2MQTT/main.py

# python3 -u main.py